# schemas/{{ cookiecutter.SCHEMA_NAME }}/crontab
#this is just an example
0 15 * * * . /home/someuser/.flockrc && cd $FLOCK_ROOT && python flock.py {{ cookiecutter.SCHEMA_NAME }}